# files/etc/s6-overlay/s6-rc.d/thttpd/run
#!/command/execlineb -P
thttpd -C /etc/thttpd/thttpd.conf
