# files/etc/s6-overlay/s6-rc.d/traefik-manager/run
#!/command/with-contenv sh
/usr/local/bin/voyager-traefik-manager --containerId=$HOSTNAME --traefikConfig=/etc/traefik/dynamic.yaml
